# postgis/cron/road_anxiety.cron
* * * * * /bin/bash /postgis/cron/road_anxiety.sh >> /var/log/cron.log 2>&1
